# nut/root/etc/services.d/nut/run
#!/usr/bin/with-contenv bash
/usr/bin/python3 /app/server.py --help