# sharepoint-downloader/root/etc/services.d/node/run
#!/usr/bin/with-contenv bash

cd /SharepointDownloader-master && /usr/bin/node .
