# lsc/root/etc/services.d/nut/run
#!/usr/bin/with-contenv bash
/app/bin/lsc -f /config -a all -n