# ssr-client/root/etc/services.d/privoxy/run
#!/usr/bin/with-contenv bash
# s6-svwait -t 10000 -u /var/run/s6/services/ssr
sleep 10000
privoxy /config/config