# ssr-client/root/etc/services.d/privoxy/run
#!/usr/bin/with-contenv bash
privoxy --no-daemon /config/config