# ssr-client/root/etc/services.d/privoxy/run
#!/usr/bin/with-contenv bash
privoxy /config/config