# ssr-client/root/etc/services.d/privoxy/run
#!/usr/bin/with-contenv bash

/usr/sbin/privoxy /config/config
