# lsc/root/etc/services.d/lsc/run
#!/usr/bin/with-contenv bash
/app/bin/lsc -f /config -s all -c all