# sasl/root/etc/services.d/nut/run
#!/usr/bin/with-contenv bash
/usr/sbin/saslauthd -a ldap -d 1