# ssr-client/root/etc/services.d/ssr/run
#!/usr/bin/with-contenv bash

/usr/bin/python /shadowsocks/local.py -c /config/config.json