# ssr-client/root/etc/services.d/ssr/run
#!/usr/bin/with-contenv bash
python /shadowsocks/local.py -c /config/config.json