# switch-lan-play/root/etc/services.d/switch-lan-play-client/run
#!/usr/bin/with-contenv bash
/usr/bin/lan-play-linux --relay-server-addr ${SERVER}:${PORT}